# shell/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; socket 文件路径
chmod=0700                       ; socket 权限

[supervisord]
logfile=/app/logs/supervisord.log ; 主日志文件
pidfile=/var/run/supervisord.pid ; pid 文件
nodaemon=true                    ; 【关键】在前台运行，防止 Docker 退出
user=root                        ; 运行用户

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; 连接 socket


[program:crew_rul]
command=bash /app/shell/start_app.sh -f /app/main.py -a crew_rul -p 37700
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/crew_rul.err.log
stdout_logfile=/app/logs/crew_rul.out.log

[program:cwru_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a cwru_cls -p 37701
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/cwru_cls.err.log
stdout_logfile=/app/logs/cwru_cls.out.log

[program:ele_metric]
command=bash /app/shell/start_app.sh -f /app/main.py -a ele_metric -p 37702
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/ele_metric.err.log
stdout_logfile=/app/logs/ele_metric.out.log

[program:industrial_meter]
command=bash /app/shell/start_app.sh -f /app/main.py -a industrial_meter -p 37703
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/industrial_meter.err.log
stdout_logfile=/app/logs/industrial_meter.out.log

[program:air_compressor_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a air_compressor_cls -p 37704
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/air_compressor_cls.err.log
stdout_logfile=/app/logs/air_compressor_cls.out.log

[program:mp_ocr]
command=bash /app/shell/start_app.sh -f /app/main.py -a mp_ocr -p 37705
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/mp_ocr.err.log
stdout_logfile=/app/logs/mp_ocr.out.log

[program:wendu_prd]
command=bash /app/shell/start_app.sh -f /app/main.py -a wendu_prd -p 37706
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/wendu_prd.err.log
stdout_logfile=/app/logs/wendu_prd.out.log

[program:zhoucheng_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a zhoucheng_cls -p 37707
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/zhoucheng_cls.err.log
stdout_logfile=/app/logs/zhoucheng_cls.out.log

[program:dianji_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a dianji_cls -p 37708
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/dianji_cls.err.log
stdout_logfile=/app/logs/dianji_cls.out.log

[program:belt_det]
command=bash /app/shell/start_app.sh -f /app/main.py -a belt_det -p 37709
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/belt_det.err.log
stdout_logfile=/app/logs/belt_det.out.log

[program:driver_belt_rul]
command=bash /app/shell/start_app.sh -f /app/main.py -a driver_belt_rul -p 37710
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/driver_belt_rul.err.log
stdout_logfile=/app/logs/driver_belt_rul.out.log

[program:chilun_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a chilun_cls -p 37711
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/chilun_cls.err.log
stdout_logfile=/app/logs/chilun_cls.out.log

[program:gangyin_ocr]
command=bash /app/shell/start_app.sh -f /app/main.py -a gangyin_ocr -p 37712
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/gangyin_ocr.err.log
stdout_logfile=/app/logs/gangyin_ocr.out.log

[program:tujiaoji_cls]
command=bash /app/shell/start_app.sh -f /app/main.py -a tujiaoji_cls -p 37713
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/tujiaoji_cls.err.log
stdout_logfile=/app/logs/tujiaoji_cls.out.log

[program:tubuji_com]
command=bash /app/shell/start_app.sh -f /app/main.py -a tubuji_com -p 37714
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/tubuji_com.err.log
stdout_logfile=/app/logs/tubuji_com.out.log

[program:duanmo_prd]
command=bash /app/shell/start_app.sh -f /app/main.py -a duanmo_prd -p 37715
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/duanmo_prd.err.log
stdout_logfile=/app/logs/duanmo_prd.out.log

[program:centrifugal_fault_det]
command=bash /app/shell/start_app.sh -f /app/main.py -a centrifugal_fault_det -p 37716
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/centrifugal_fault_det.err.log
stdout_logfile=/app/logs/centrifugal_fault_det.out.log

[program:fan_noise_flask]
command=bash /app/shell/start_app.sh -f /app/main.py -a fan_noise_flask -p 37717
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/fan_noise_flask.err.log
stdout_logfile=/app/logs/fan_noise_flask.out.log

[program:jiaobanji_prd]
command=bash /app/shell/start_app.sh -f /app/main.py -a jiaobanji_prd -p 37718
directory=/app
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/app/logs/jiaobanji_prd.err.log
stdout_logfile=/app/logs/jiaobanji_prd.out.log
